# etc/profile.d/vdpau_driver.sh
# Local script to set VDPAU_DRIVER to va_gl for Intel Graphics Cards
export VDPAU_DRIVER=va_gl
