# usr/lib/systemd/system/hwclocksync.timer
[Unit]
Description=Sync hardware clock to system clock

[Timer]
OnBootSec=15m
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
