# usr/lib/systemd/system/hwclocksync.service
[Unit]
Description=Sync hardware clock to system clock

[Service]
Type=oneshot
ExecStart=/usr/bin/hwclock --systohc
IOSchedulingClass=idle
